# quickcall.rules
SUBSYSTEM=="usb", ACTION=="add", ATTR{idVendor}=="046d", ATTR{idProduct}=="08d5", RUN+="@QUICKCALLD@ -u -v %S%p"
